# deploy/com.andres.daily-news.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.andres.daily-news</string>

    <!-- Run the daily routine: generate today's edition, commit the new snapshot
         (the system's state), then push so GitHub Pages rebuilds. The wrapper sets
         PATH itself, but we also set it below since launchd does not source your
         shell profile. PATH must include the dir holding `bun` and the `claude`
         CLI (Homebrew: /opt/homebrew/bin) because research.ts spawns claude. -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/andres/Projects/daily-news/deploy/daily-run.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/andres/Projects/daily-news</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/local/sbin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/andres/.bun/bin</string>
    </dict>

    <!-- No calendar trigger here ON PURPOSE. The 07:00 schedule lives in the
         Claude Code routine "daily-news" (visible in the Routines panel), which
         fires `launchctl start com.andres.daily-news`. This job is the detached,
         robust EXECUTOR (the ~15-20 min generation outlives any app session and
         isn't bound by an agent's command timeout). Keeping the schedule in only
         one place avoids double-runs. To go back to a pure-launchd 07:00 trigger,
         re-add StartCalendarInterval { Hour 7, Minute 0 } and pause the routine. -->

    <!-- Do NOT run on load: only run when explicitly triggered (by the routine,
         or manually via `launchctl start com.andres.daily-news`). -->
    <key>RunAtLoad</key>
    <false/>

    <key>StandardInputPath</key>
    <string>/dev/null</string>

    <key>StandardOutPath</key>
    <string>/Users/andres/Library/Logs/com.andres.daily-news.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/andres/Library/Logs/com.andres.daily-news.err.log</string>
</dict>
</plist>
